# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command= docker-php-entrypoint php-fpm
autostart=true
autorestart=true

[program:npm]
command= npm run dev
